# resources/fty-alert-stats.service.in
[Unit]
Description=@PROJECT_NAME@ service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

After=fty-asset.service
Requires=fty-asset.service
BindsTo=fty-asset.service

[Service]
Type=simple
User=@AGENT_USER@
Restart=always
Environment='SYSTEMD_UNIT_FULLNAME=%n'
Environment="prefix=/usr"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf

# etc
ExecStartPre=+/bin/chmod 755 @AGENT_ETC_DIR@
ExecStartPre=+/bin/chown -R www-data @AGENT_ETC_DIR@/*
ExecStartPre=+/bin/chmod 644 @AGENT_ETC_DIR@/*
# var/lib/fty
#ExecStartPre=+/bin/chown -R @AGENT_USER@ @AGENT_VAR_DIR@
#ExecStartPre=+/bin/chmod 771 @AGENT_VAR_DIR@

# exec
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@ -c @AGENT_ETC_DIR@/@PROJECT_NAME@.cfg

[Install]
WantedBy=bios.target
